# tests/crontab
* * * * * echo "Test every minute" >> /home/test/minute.log
