# tests/crontab
# minute hour dayOfMonth monthOfYear (1=mon; 7=sun) dayOfWeek
*	*	*	*	*	echo "Test every minute" >> /home/test/minute.log
